# systemd/local-ai-voice-tts-chatterbox.service
[Unit]
Description=Local AI Voice TTS worker - Chatterbox
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=local-ai-voice
Group=local-ai-voice
SupplementaryGroups=video render
WorkingDirectory=/opt/local-ai-voice/app/workers/tts-chatterbox
EnvironmentFile=/opt/local-ai-voice/config/local-ai-voice.env
Environment=TTS_PROVIDER=chatterbox
Environment=TTS_CHATTERBOX_BIND_HOST=127.0.0.1
Environment=TTS_CHATTERBOX_BIND_PORT=8001
Environment=TTS_BIND_HOST=127.0.0.1
Environment=TTS_BIND_PORT=8001
ExecStart=/opt/local-ai-voice/workers/tts/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001 --log-level info
Restart=on-failure
RestartSec=5
TimeoutStopSec=60
LimitNOFILE=65536
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/local-ai-voice

[Install]
WantedBy=multi-user.target
